# roles/systemd/files/suspend@.service
[Unit]
Description=Lock the screen
Before=sleep.target

[Service]
User=%I
Type=forking
# https://github.com/google/xsecurelock#known-compatibility-issues
Environment="DISPLAY=:0" "XSECURELOCK_COMPOSITE_OBSCURER=0"
ExecStart=xsecurelock

[Install]
WantedBy=sleep.target